# configs/demo-app-periodic.service
[Unit]
Description=Periodic Demo App

[Service]
Type=oneshot
ExecStart=/usr/local/bin/log-demo.sh
